# doorbell.service
# ═══════════════════════════════════════════════════════════════════════════
# SECTION 1: [Unit] - Identity and Dependencies
# ═══════════════════════════════════════════════════════════════════════════
[Unit]

Description=Museum Exhibit Camera Button
# ┗━━ Human-readable name. Shows up in:
#     • systemctl status
#     • journalctl logs
#     • Error messages
#     Think of it like a comment, but one the system actually uses.

After=multi-user.target
# ┗━━ "Don't start me until multi-user.target is reached"
#     
#     Targets are like milestones in the boot process:
#     
#     Boot start ──► basic.target ──► multi-user.target ──► graphical.target
#                    (minimal)        (full system,         (desktop)
#                                      no desktop)
#     
#     For a headless pi setup, multi-user.target means:
#     • Filesystem is mounted
#     • Network is configured
#     • SSH is available
#     • System is ready for normal operation
#     
#     Your service starts AFTER this milestone.
#     
#     Some common directive types:
#     • Wants=    Soft: "start that too, but I'll survive without it"
#     • Requires= Hard: "if that fails, I fail too"
#     • After=    Order: "start me after that" (but doesn't pull it in)

StartLimitIntervalSec=0
# ┗━━ "Never give up restarting this service."
#
#     systemd has a safety valve: by default, if a service restarts too many
#     times too quickly (5 times within 10s) it assumes the situation is
#     hopeless, gives up, and leaves the service DEAD. For a normal desktop app
#     that's sensible. For an unattended 24/7 exhibit it is the WRONG behavior:
#     a camera that is simply slow to enumerate at boot could trip this limit,
#     and the exhibit would then stay dark until a human physically intervenes.
#
#     Setting this to 0 disables the rate limit entirely. Combined with
#     Restart=always (below), systemd will keep retrying FOREVER, every
#     RestartSec seconds, no matter how many times it has failed. The exhibit
#     always claws its way back on its own.


# ═══════════════════════════════════════════════════════════════════════════
# SECTION 2: [Service] - How to Run
# ═══════════════════════════════════════════════════════════════════════════
[Service]

Type=notify
# ┗━━ service category/behavior specified here
#
#     Service types:
#     • simple:  ExecStart (see below) is the service (most Python scripts)
#     • forking: ExecStart spawns a child and exits (traditional daemons)
#     • oneshot: Runs once and exits (setup scripts)
#     • notify:  Like simple, but the service tells systemd when it's ready
#                (and can ping a watchdog) via sd_notify messages
#
#     We use 'notify' (not 'simple') because the script talks back to systemd:
#       • it sends READY=1 once the camera is initialized, so "active" in
#         systemctl status genuinely means the exhibit is up
#       • it sends WATCHDOG=1 every loop so systemd can detect a hang
#         (see WatchdogSec below)
#     This requires the python3-systemd package (see README install step).

WatchdogSec=15
# ┗━━ "If I don't hear a WATCHDOG=1 ping for 15 seconds, assume the process is
#     hung and restart it."
#
#     Restart=always (below) already handles a process that CRASHES (exits).
#     But a worse failure for an exhibit is a process that is still alive yet
#     FROZEN -- e.g. a camera/DRM call that blocks forever. systemd can't see
#     that from the outside; the process looks fine.
#
#     With this set, the script must check in at least every 15s by sending
#     WATCHDOG=1 from its main loop. The loop cycles every 0.1s, so in healthy
#     operation it pings ~150x per window -- huge margin. The moment the loop
#     stops cycling (a hang), the pings stop, 15s elapses, and systemd kills +
#     restarts the service automatically. Self-healing with no human needed.
#
#     NOTE: the watchdog clock only starts after the script sends READY=1, so a
#     slow camera at boot won't trip it; startup is governed by TimeoutStartSec
#     (default 90s), which the camera retry loop comfortably fits inside.

User=doorbell
Group=doorbell
SupplementaryGroups=video render gpio
# ┗━━ Run as the unprivileged 'doorbell' user instead of root.
#
#     By default (no User= line) a service runs as ROOT — full control of the
#     entire machine. If the script or a library it pulls in ever misbehaves or
#     is compromised, root means there are no guardrails. Professional practice
#     is to run with the least privilege the job actually needs.
#
#     This camera app needs exactly three capabilities beyond a normal user:
#       • video, render → talk to the camera and the DRM/KMS display pipeline
#                         that draws the live preview on the monitor
#       • gpio          → read the arcade button wired to GPIO 17
#
#     SupplementaryGroups grants those specific groups and nothing more.
#     PREREQUISITE: the 'doorbell' user must exist and own the script directory
#     (it does, since the script lives under /home/doorbell/).

ExecStart=/usr/bin/python3 /home/doorbell/doorbell_camera/doorbell_camera.py
# ┗━━ The command to run. MUST use absolute paths. Make sure second directory 
#     path matches path of python script to run exactly
#     
#     Equivalent to typing this in a terminal:
#     $ /usr/bin/python3 /home/doorbell/doorbell_camera/doorbell_camera.py

WorkingDirectory=/home/doorbell/doorbell_camera
# ┗━━ "cd to this directory before running ExecStart"
#     
#     Affects:
#     • Where relative file paths resolve
#     • Where Python looks for modules (if you import local files)
#     
#     Good practice even if your script uses absolute paths.

# ─────────────────────────────────────────────────────────────────────────
# Restart behavior (CRITICAL for exhibits)
# ─────────────────────────────────────────────────────────────────────────

Restart=always
# ┗━━ "If the process exits FOR ANY REASON, start it again"
#     
#     Restart policies:
#     • no:         Never restart (default)
#     • always:     Always restart 
#     • on-failure: Only restart if exit code != 0
#     • on-abnormal: Restart on signal/timeout/watchdog
#     
#     'always' means:
#     • Crash? Restart.
#     • Clean exit? Restart.
#     • Killed by signal? Restart.
#     • Only way to stop: 'systemctl stop'
#     

RestartSec=2
# ┗━━ "Wait 2 seconds before restarting"

# ─────────────────────────────────────────────────────────────────────────
# Clean shutdown
# ─────────────────────────────────────────────────────────────────────────

KillSignal=SIGTERM
# ┗━━ "When I need to stop this service, send SIGTERM first"
#     
#     SIGTERM = "Please shut down gracefully"
#     SIGKILL = "Die immediately" (cannot be caught)
#     
#     doorbell_camera.py script has:
#         signal.signal(signal.SIGTERM, handle_shutdown)
#     
#     So when systemd sends SIGTERM, your script:
#     1. Sets shutdown_requested = True
#     2. Exits the main loop cleanly
#     3. Runs the finally: block
#     4. Closes camera resources
#     5. Exits
#     
#     (handshake between systemd and your code)

TimeoutStopSec=5
# ┗━━ "Give the service 5 seconds to shut down gracefully"
#     
#     Timeline when you run 'systemctl stop':
#     
#     0s: systemd sends SIGTERM
#     │
#     │   Your script is shutting down...
#     │   Closing camera...
#     │   Cleanup...
#     │
#     5s: If still running, systemd sends SIGKILL (forced death)
#     
#     5 seconds is plenty for your script. Camera cleanup is fast.

# ─────────────────────────────────────────────────────────────────────────
# Logging
# ─────────────────────────────────────────────────────────────────────────

StandardOutput=journal
StandardError=journal
# ┗━━ "Send all print() output to the systemd journal"
#     
#     The journal is systemd's logging system. All your print()
#     statements become queryable log entries.
#     
#     View logs:
#         journalctl -u doorbell.service        # All logs
#         journalctl -u doorbell.service -f     # Follow (like tail -f)
#         journalctl -u doorbell.service --since "1 hour ago"
#     
#     Without this, output would go nowhere (no terminal to print to).

# ─────────────────────────────────────────────────────────────────────────
# Environment sanity
# ─────────────────────────────────────────────────────────────────────────

Environment=PYTHONUNBUFFERED=1
# ┗━━ "Set this environment variable for the process"
#     
#     Python normally buffers stdout for performance. This means
#     print() statements might not appear in logs immediately.
#     
#     PYTHONUNBUFFERED=1 disables this buffering, so:
#     • Logs appear in real-time
#     • If the script crashes, you see the last message
#     • Debugging is much easier
#     
#     ALWAYS use this for systemd Python services.


# ═══════════════════════════════════════════════════════════════════════════
# SECTION 3: [Install] - Boot Integration
# ═══════════════════════════════════════════════════════════════════════════
[Install]

WantedBy=multi-user.target
# ┗━━ "When enabled, attach me to multi-user.target"
#     
#     This section is ONLY used by:
#         systemctl enable doorbell.service
#         systemctl disable doorbell.service
#     
#     When you run 'enable', systemd creates a symlink:
#         /etc/systemd/system/multi-user.target.wants/doorbell.service
#         → /etc/systemd/system/doorbell.service
#     
#     This means: "When booting to multi-user.target, also start
#     doorbell.service"
#     
#     For a headless Pi, multi-user.target is the right choice.
#     (graphical.target would be for desktop systems)
